# configs/crontab
# comcast data usage
0,30 * * * *  python3 /scripts/comcast.py